# com.jashan.yield-finder-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jashan.yield-finder-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/start-server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPO_PATH</string>

  <!-- Start on login and restart if it crashes -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Redirect output to log files -->
  <key>StandardOutPath</key>
  <string>REPO_PATH/server.log</string>
  <key>StandardErrorPath</key>
  <string>REPO_PATH/server.log</string>

  <!-- Wait 5s before restarting on crash -->
  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
